# kindle-tools/ssh/koreader-ssh-autostart.conf
# Key-only SSH at boot, using KOReader's patched Dropbear and authorized_keys.
# The wait covers boots where /mnt/us is temporarily exported over USB.

start on framework_ready
stop on stopping system

respawn
respawn limit 5 60

pre-start script
    attempts=0
    while [ "$attempts" -lt 300 ]; do
        if [ -x /mnt/us/koreader/dropbear ] && [ -s /mnt/us/koreader/settings/SSH/authorized_keys ]; then
            break
        fi
        attempts=$((attempts + 1))
        sleep 1
    done

    [ -x /mnt/us/koreader/dropbear ] || exit 1
    [ -s /mnt/us/koreader/settings/SSH/authorized_keys ] || exit 1

    iptables -D INPUT -p tcp --dport 2222 -m conntrack --ctstate NEW,ESTABLISHED -j ACCEPT >/dev/null 2>&1 || true
    iptables -D OUTPUT -p tcp --sport 2222 -m conntrack --ctstate ESTABLISHED -j ACCEPT >/dev/null 2>&1 || true
    iptables -A INPUT -p tcp --dport 2222 -m conntrack --ctstate NEW,ESTABLISHED -j ACCEPT
    iptables -A OUTPUT -p tcp --sport 2222 -m conntrack --ctstate ESTABLISHED -j ACCEPT
end script

script
    cd /mnt/us/koreader
    exec ./dropbear -F -E -R -p2222 -P /var/run/koreader-ssh-boot.pid -s
end script

post-stop script
    iptables -D INPUT -p tcp --dport 2222 -m conntrack --ctstate NEW,ESTABLISHED -j ACCEPT >/dev/null 2>&1 || true
    iptables -D OUTPUT -p tcp --sport 2222 -m conntrack --ctstate ESTABLISHED -j ACCEPT >/dev/null 2>&1 || true
    rm -f /var/run/koreader-ssh-boot.pid
end script
