# kindle-tools/koreader-personalization/koreader-autostart.conf
# KOReader autostart for Kindle Paperwhite 3 / firmware 5.16.2.1.1.
# The job remains inert unless /mnt/us/ENABLE_KOREADER_AUTOSTART exists.

start on framework_ready and started lab126_gui and started kppmainapp
stop on stopping system

task

script
    # During USB mass-storage mode /mnt/us belongs to the computer. Wait for it
    # to return instead of silently giving up during a cable-connected boot.
    attempts=0
    while [ "$attempts" -lt 300 ]; do
        if [ -f /mnt/us/ENABLE_KOREADER_AUTOSTART ] && [ -f /mnt/us/koreader/koreader.sh ]; then
            if [ -x /usr/sbin/eips ] && [ -f /mnt/us/koreader/screensavers/library-mountain-pw3.png ]; then
                /usr/sbin/eips -g /mnt/us/koreader/screensavers/library-mountain-pw3.png
                sleep 2
            fi
            exec /bin/sh /mnt/us/koreader/koreader.sh --kual
        fi
        attempts=$((attempts + 1))
        sleep 1
    done
end script
